# chef/cookbooks/app/templates/default/etc/profile.d/php-aliases.sh.erb
alias php-debug='XDEBUG_CONFIG="" PHP_IDE_CONFIG="serverName=<%= @server_name_in_ide %>" php'
alias php-profile='XDEBUG_CONFIG="profiler_enable=1" PHP_IDE_CONFIG="serverName=<%= @server_name_in_ide %>" php'